# launchd/com.sentinel.scanner.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sentinel.scanner.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__SCANNER_DIR__/build/scanner</string>
        <string>monitor</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__SCANNER_DIR__/scanner-data/scanner.log</string>

    <key>StandardErrorPath</key>
    <string>__SCANNER_DIR__/scanner-data/scanner.log</string>
</dict>
</plist>
